# sandbox/aio/docker/context/rootfs/opt/gem/supervisord/agent-browser.conf
[program:agent-browser]
environment=HOME="/home/%(ENV_USER)s",USER="%(ENV_USER)s",PATH="/usr/local/bin:/usr/bin:/bin:/home/%(ENV_USER)s/.npm-global/bin"
command=/opt/gem/agent-browser-init.sh
autostart=%(ENV_AUTOSTART_BROWSER)s
autorestart=false
startsecs=0
exitcodes=0
priority=200
user=%(ENV_USER)s
stdout_logfile=%(ENV_LOG_DIR)s/agent-browser.log
stdout_logfile_maxbytes=5MB
stdout_logfile_backups=2
redirect_stderr=true
