# sandbox/aio/docker/context/rootfs/opt/gem/supervisord/supervisord.nodejs_repl_24.conf
[program:nodejs-repl-server-24]
environment=HOME="/home/%(ENV_USER)s",USER="%(ENV_USER)s",NODEJS_REPL_PORT="%(ENV_NODEJS_REPL_PORT_24)s"
command=/usr/local/bin/node24 /opt/repl-servers/nodejs/server.js
; WORKSPACE must be exported by run.sh before supervisord starts (no fallback syntax in supervisord).
directory=%(ENV_WORKSPACE)s
autostart=%(ENV_AUTOSTART_NODEJS_REPL_24)s
autorestart=true
priority=90
user=%(ENV_USER)s
stdout_logfile=%(ENV_LOG_DIR)s/nodejs-repl-server-24.log
stdout_logfile_maxbytes=10MB
stdout_logfile_backups=10
redirect_stderr=true
